# src/etc/runit/one-contextd/run
#!/bin/sh
exec 2>&1
/usr/local/sbin/one-contextd all reconfigure
exec chpst -b one-contextd pause
